# Tools/tg-bot/v1.1.1/telegram-bot.service
[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app/bot
ExecStart=/usr/bin/python3 /app/bot/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
